# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/php-fpm.err.log
stdout_logfile=/var/log/supervisor/php-fpm.out.log

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log

[program:laravel-setup]
command=/bin/sh -c "sleep 15 && echo 'Testing database connection...' && php artisan migrate:status || echo 'Database not ready, skipping migrations' && php artisan config:cache && php artisan route:cache && php artisan view:cache || echo 'Cache commands failed'"
autostart=true
autorestart=false
stderr_logfile=/var/log/supervisor/laravel-setup.err.log
stdout_logfile=/var/log/supervisor/laravel-setup.out.log
